# install/background_jobs/supervisor_conf/common-friends.conf
[program:common-friends]
command=php /var/www/localhost/htdocs/artisan command:common-friends
autostart=true
autorestart=true
stdout_logfile=/tmp/common-friends.log
